# monitoring/systemd/openjarvis.service
[Unit]
Description=OpenJarvis AI Agent
After=network.target

[Service]
Type=simple
User=edvisingu
WorkingDirectory=/opt/edvisingu/openjarvis
ExecStart=/opt/edvisingu/openjarvis/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
